# c2/files/master.conf.erb
start on filesystem and net-device-up IFACE=eth0

pre-start script

bash << "EOF"
  mkdir -p <%= log_root %>
  chown -R <%= user %> <%= log_root %>
EOF

end script