# python/files/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn

script
  sudo -u <%= user %> -i
  cd <%= engine.directory %>
  export PORT=<%= port %>;<% engine.environment.each_pair do |var,env| %> export <%= var.upcase %>='<%= env %>'; <% end %>
  . venv/bin/activate
  exec <%= process.command %> >> <%= log_root %>/<%=process.name%>-<%=num%>.log 2>&1
end script