# custom/files/process.conf.erb
start on starting <%= app %>-<%= process.name %>
stop on stopping <%= app %>-<%= process.name %>
respawn

setuid <%= user %>
setgid <%= user %>

script
  . .opdemand
  export PATH=$PATH:<%= engine.directory %>/bin
  cd <%= engine.directory %>
  exec <%= process.command %> >> <%= log_root %>/<%=process.name%>-<%=num%>.log 2>&1
end script